# debian/prerm
#!/bin/bash
set -e

# Disable and stop services on removal
if [ -d /run/systemd/system ]; then
    for service in cloud-init-local cloud-init cloud-config cloud-final; do
        systemctl --no-reload disable "${service}.service" >/dev/null 2>&1 || true
        systemctl stop "${service}.service" >/dev/null 2>&1 || true
    done
fi

exit 0
